# projects/quiz-video-generator/deployment/quiz-generator.service
[Unit]
Description=Quiz Video Generator Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/quiz-generator
Environment="PATH=/var/www/quiz-generator/venv/bin"
ExecStart=/var/www/quiz-generator/venv/bin/python3 app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
